# happychina-yiimp/docker/yiimp/config/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid
minfds=4096

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:memcached]
command=/usr/bin/memcached -m 64 -p 11211 -u memcache -l 127.0.0.1
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:apache]
command=/usr/sbin/apachectl -D FOREGROUND
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:yiimp-main]
command=/bin/bash /var/www/main.sh
directory=/var/www
user=www-data
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:yiimp-loop2]
command=/bin/bash /var/www/loop2.sh
directory=/var/www
user=www-data
autostart=true
autorestart=true
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:yiimp-blocks]
command=/bin/bash /var/www/blocks.sh
directory=/var/www
user=www-data
autostart=true
autorestart=true
priority=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:stratum-3331]
command=/bin/bash -lc 'ulimit -c 0; exec /usr/local/bin/stratum /etc/yiimp/stratum/scrypt-3331.conf'
autostart=true
autorestart=true
priority=100
startsecs=2
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:stratum-3332]
command=/bin/bash -lc 'ulimit -c 0; exec /usr/local/bin/stratum /etc/yiimp/stratum/scrypt-3332.conf'
autostart=true
autorestart=true
priority=101
startsecs=2
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:stratum-3333]
command=/bin/bash -lc 'ulimit -c 0; exec /usr/local/bin/stratum /etc/yiimp/stratum/scrypt-3333.conf'
autostart=true
autorestart=true
priority=102
startsecs=2
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:stratum-3334]
command=/bin/bash -lc 'ulimit -c 0; exec /usr/local/bin/stratum /etc/yiimp/stratum/scrypt-3334.conf'
autostart=true
autorestart=true
priority=103
startsecs=2
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:stratum-3335]
command=/bin/bash -lc 'ulimit -c 0; exec /usr/local/bin/stratum /etc/yiimp/stratum/scrypt-3335.conf'
autostart=true
autorestart=true
priority=104
startsecs=2
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:stratum-3336]
command=/bin/bash -lc 'ulimit -c 0; exec /usr/local/bin/stratum /etc/yiimp/stratum/scrypt-3336.conf'
autostart=true
autorestart=true
priority=105
startsecs=2
startretries=999
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
